# knowledge/scripts/com.librarian.curate-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.librarian.curate-server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/mcginleyr1/develop/the-librarian/knowledge/scripts/curate-server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mcginleyr1/develop/the-librarian</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mcginleyr1/develop/the-librarian/knowledge/curate-server.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mcginleyr1/develop/the-librarian/knowledge/curate-server.log</string>
</dict>
</plist>
